# filesync.service
[Unit]
Description=Run rsync backup every 5 minutes
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncscript
